# deploy/systemd/chief-agent.service
[Unit]
Description=Chief of Staff Agent Server
After=litellm.service

[Service]
User=nonroot
WorkingDirectory=/app/workspace/agents/server
Environment=LITELLM_URL=http://127.0.0.1:4000
Environment=LITELLM_KEY=sk-cos-local-dev
Environment=OPENAI_API_KEY=sk-cos-local-dev
ExecStart=/app/venv/bin/uvicorn server:app \
  --host 127.0.0.1 \
  --port 8002
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
