# deploy/systemd/headroom.service
[Unit]
Description=Headroom Context Compression Proxy
After=network.target

[Service]
User=nonroot
Environment=HEADROOM_TELEMETRY=off
Environment=OPENAI_TARGET_API_URL=http://127.0.0.1:8080/v1
ExecStart=/app/venv/bin/headroom proxy \
  --host 127.0.0.1 \
  --port 8787 \
  --no-telemetry
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
